# systemd/assistantd.service
[Unit]
Description=Local AI Assistant Daemon (Scaffold)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/var/lib/local-ai-assistant
EnvironmentFile=/etc/local-ai-assistant/assistantd.env
ExecStart=/usr/local/bin/assistantd --config /etc/local-ai-assistant/assistantd.env
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# TODO(implementation-playbook):
# - Confirm least-privilege user/group and runtime directory ownership model.
# - Add explicit hardening options (NoNewPrivileges, ProtectSystem, PrivateTmp, etc.).
# - Decide whether watchdog notifications are needed for process health.

[Install]
WantedBy=multi-user.target
